# crontab
*/10 * * * * echo "[$(/home/wbousfir/born2beroot/monitoring.sh)]" |wall
